# ansible/templates/etc/network/if-up.d/vxlan_routing.j2
#!/bin/sh

{% for vxlan in network_configuration.routers.value[inventory_hostname].vxlans %}
route add -net {{ vxlan.remote_subnet }} gw {{ vxlan.next_hop }} dev {{ br_vxlan }} || true
{% endfor %}